# minidlna.conf
[program:start]
stdout_logfile = /config/supervisord-start-stdout.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 5
redirect_stderr = true
user = nobody
command = /run/minidlna/start.sh