# build/minidlna.conf
[program:crond]
user = root
autorestart = False
command = /root/crond.sh
umask = 000

[program:start]
user = nobody
startsecs = 0
autorestart = False
command = /home/nobody/minidlna.sh
umask = 000



