# build/minidlna.conf
[program:crond]
user = root
autorestart = False
command = /root/crond.sh
umask = 000

[program:start]
user = nobody
startsecs = 0
autorestart = False
command = /home/nobody/minidlna.sh
umask = 000

[program:shutdown-script]
autorestart = false
startsecs = 0
user = root
command = /usr/local/bin/shutdown.sh 'minidlnad'
umask = 000

