# setup/minidlna.conf
[program:start]
user = root
startsecs = 0
autorestart = False
command = /home/nobody/setup.sh

[program:minidlna]
user = nobody
command = /usr/bin/minidlnad -S -R -f /config/minidlna.conf